# script/routine/crontab
* 01 * * 02,05 sh /g6/liqq/lzhenn/routine.sh >/dev/null 2>&1
